# docker/data_scraper/crontab
0 18 * * 1-5 root cd /finance && make scrape scraper=cboe && make scrape scraper=tiingo
0 0 1 * * root cd /finance && make aggregate && make backup

