# docker/data_scraper/crontab
0 18 * * 1-5 root cd /finance && make scrape scraper=cboe && make scrape scraper=tiingo && make backup

